# tests/ini2dict/examples/file.ini
[program:github_push]
command = python -m my_daemons.github.webhooks.push
autorestart = true
startretries = 10
